# etc/sv/httpd/run
#!/bin/sh

exec /usr/sbin/httpd -DFOREGROUND 2>&1
